# templates/rc.local.erb
#!/bin/sh -e

iptables -t nat -F PREROUTING
iptables -t nat -P PREROUTING ACCEPT

iptables -t filter -F FORWARD
iptables -t filter -P FORWARD ACCEPT

iptables -t mangle -F
iptables -t nat -F POSTROUTING

# Mark all packets that come from a vpn (input) interface so we can SNAT them
# later in the POSTROUTING of the nat table.
<% @vpn_devs.each do |dev| %>
  iptables -t mangle -A PREROUTING -i <%= dev %> -j MARK --set-mark 42
<% end %>

# SNAT all packets that are marked. Marked packets come from one of the vpn
# intput interfaces.
iptables -t nat -A POSTROUTING \
  -m mark --mark 42 \
  -j SNAT --to <%= @inet_add %>.<%= @inet_min %>-<%= @inet_add %>.<%= @inet_max %>

# Do not route to destinations that are part of RFC1918:
# https://tools.ietf.org/html/rfc1918
<% @private_networks.each do |network| %>
  iptables -t filter -A FORWARD -o <%= @inet_dev %> -d <%= network %> \
    -j REJECT --reject-with icmp-net-prohibited
<% end %>

# Circumvent Path MTU Discovery issues with MSS Clamping:
# http://lartc.org/howto/lartc.cookbook.mtu-mss.html
iptables -t filter -A FORWARD -o <%= @inet_dev %> \
  -p tcp --tcp-flags SYN,RST SYN -j TCPMSS --clamp-mss-to-pmtu

ip rule add iif <%= @inet_dev %> lookup reverseroutes
ip rule add iif lo lookup reverseroutes

# Increasing the size of nf_conntrack hash-table so we don't drop packets
# because of a full table.
echo 32768 > /sys/module/nf_conntrack/parameters/hashsize

# ensure sysctl settings are set on boot
# https://bugs.launchpad.net/ubuntu/+source/procps/+bug/50093

/etc/init.d/procps restart

exit 0
